# .github/server-release/orynivo-server-updater.path
[Unit]
Description=Watch for verified Orynivo Server update requests

[Path]
PathExists=/var/lib/orynivo-server/updates/ready
Unit=orynivo-server-updater.service

[Install]
WantedBy=multi-user.target
