# .github/server-release/orynivo-server.service
[Unit]
Description=Orynivo Music Server
Documentation=https://github.com/bschlaack/Orynivo
After=network.target

[Service]
Type=simple
User=orynivo-server
Group=orynivo-server
WorkingDirectory=/usr/lib/orynivo-server
ExecStart=/usr/lib/orynivo-server/Orynivo.Server
# systemd creates /var/lib/orynivo-server owned by the service user; the server
# stores its database, caches, and downloaded artwork there. Without a writable
# data directory the --no-create-home service user cannot create the SQLite
# database under $HOME and the process aborts on startup.
StateDirectory=orynivo-server
Environment=ORYNIVO_DATA_DIR=/var/lib/orynivo-server
EnvironmentFile=-/etc/default/orynivo-server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orynivo-server

[Install]
WantedBy=multi-user.target
